# supervisor/conf.d/celery.conf
[program:celery]
command=pipenv run celery -A opps.tasks.celery  worker -l info
directory=/usr/local/2opps/
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
