# supervisor/conf.d/opps.conf
[program:opps]
command=pipenv run gunicorn -w 1 -b 0.0.0.0:5000 wsgi:app
directory=/usr/local/2opps/
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
