# supervisor/conf.d/rabbitmq3.conf
[program:rabbitmq3]
command=rabbitmq-server
directory=/opt/pro/rabbitmq3/sbin/
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
